# projects/project-popcorn/config/com.popcorn.radar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- 서비스 식별자 -->
    <key>Label</key>
    <string>com.popcorn.radar</string>

    <!-- 실행할 프로그램 -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/sokchokim/sokcho_code/MLAF/projects/project-popcorn/.venv/bin/python</string>
        <string>/Users/sokchokim/sokcho_code/MLAF/projects/project-popcorn/run_daily.py</string>
        <string>--ministry</string>
        <string>산업통상부</string>
    </array>

    <!-- 작업 디렉토리 -->
    <key>WorkingDirectory</key>
    <string>/Users/sokchokim/sokcho_code/MLAF/projects/project-popcorn</string>

    <!-- 실행 시간: 매일 09:00 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- 환경변수 (선택) -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- 표준 출력 로그 -->
    <key>StandardOutPath</key>
    <string>/Users/sokchokim/sokcho_code/MLAF/projects/project-popcorn/logs/launchd_stdout.log</string>

    <!-- 표준 에러 로그 -->
    <key>StandardErrorPath</key>
    <string>/Users/sokchokim/sokcho_code/MLAF/projects/project-popcorn/logs/launchd_stderr.log</string>

    <!-- 실패 시 재시작 안 함 -->
    <key>KeepAlive</key>
    <false/>

    <!-- 부팅 시 자동 실행 안 함 -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
